# celery.service
[Unit]
Description=Celery Worker for LootLink
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/lootlink
Environment="PATH=/opt/lootlink/venv/bin"
ExecStart=/opt/lootlink/venv/bin/celery -A config worker --loglevel=info --logfile=/var/log/lootlink/celery-worker.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

